# deployment/model-server.service
[Unit]
Description=Model Server - Generic Flask application server
After=network.target

[Service]
Type=simple
User=modelserver
Group=modelserver
WorkingDirectory=/opt/stats.server
Environment="PATH=/opt/stats.server/venv/bin"
Environment="PORT=8085"
EnvironmentFile=/opt/stats.server/.env
ExecStart=/opt/stats.server/venv/bin/gunicorn \
    --bind 0.0.0.0:${PORT} \
    --workers ${WORKERS} \
    --threads ${THREADS} \
    --timeout ${TIMEOUT} \
    --access-logfile /var/log/stats.server/access.log \
    --error-logfile /var/log/stats.server/error.log \
    --log-level ${GUNICORN_LOG_LEVEL} \
    --chdir /opt/stats.server/src \
    main:app

# Restart policy
Restart=always
RestartSec=10s

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/stats.server

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
